# ironyard.service
[Unit]
Description=Publication Monitoring Service

[Service]
Type=simple
ExecStart=/home/bot/run.sh

[Install]
WantedBy=multi-user.target
